# services/network.service
[Service]
Name=Network
Exec=/usr/bin/sleep 2
Log=/tmp/network.log
